# lfc-agent.service
[Unit]
Description=Sports Team Daily Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/sports_agent
EnvironmentFile=/opt/sports_agent/.env
ExecStart=/opt/sports_agent/venv/bin/python src/sports_agent.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
